# scripts/exporters/gpu-exporter.service.template
[Unit]
Description=NVIDIA GPU Metrics Exporter for Prometheus
Documentation=https://github.com/xencon/aixcl
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 @INSTALL_PATH@/scripts/exporters/gpu-exporter.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
